# router/files/etc/uci-defaults/81-network
set -x

. /root/.env

if [ -n "$OP_LAN_IP_PREFIX" ]; then
  config_load network
  uci set network.lan.ipaddr="${OP_LAN_IP_PREFIX}.1.1"
  uci commit network
fi

######################################
uci batch << EOI
set network.switch_vlan_ap.description='LAN and home wifi'
set network.@switch_vlan[0].vlan='1'
set network.@switch_vlan[0].vid='1'
set network.@switch_vlan[0].ports='3 4 6t'

set network.@switch_vlan[1].vlan='2'
set network.@switch_vlan[1].vid='2'
set network.@switch_vlan[1].ports='0t 5'

set network.switch_vlan_guest.description='guest'
set network.switch_vlan_guest=switch_vlan
set network.switch_vlan_guest.device='switch0'
set network.switch_vlan_guest.vlan='3'
set network.switch_vlan_guest.vid='3'
set network.switch_vlan_guest.ports='6t'

set network.switch_vlan_iot.description='iot'
set network.switch_vlan_iot=switch_vlan
set network.switch_vlan_iot.device='switch0'
set network.switch_vlan_iot.vlan='4'
set network.switch_vlan_iot.vid='4'
set network.switch_vlan_iot.ports='1 2 6t'

commit network

set network.br_guest=device
set network.br_guest.type='bridge'
set network.br_guest.name='br-guest'
set network.br_guest.ports='eth1.3'
commit network

set network.br_iot=device
set network.br_iot.type='bridge'
set network.br_iot.name='br-iot'
set network.br_iot.ports='eth1.4'
commit network

set network.lan.device='br-lan'
commit network
EOI

######################################
uci batch << EOI
set firewall.zone_iot=zone
set firewall.zone_iot.name='iot'
set firewall.zone_iot.input='ACCEPT'
set firewall.zone_iot.output='ACCEPT'
set firewall.zone_iot.forward='REJECT'
set firewall.forwarding_iot=forwarding
set firewall.forwarding_iot.src='iot'
set firewall.forwarding_iot.dest='wan'
commit firewall

set dhcp.iot=dhcp
set dhcp.iot.interface='iot'
set dhcp.iot.start='100'
set dhcp.iot.limit='150'
set dhcp.iot.leasetime='12h'
commit dhcp

add_list firewall.zone_iot.network='iot'
commit firewall

set network.iot=interface
set network.iot.proto='static'
set network.iot.device='br-iot'
set network.iot.ipaddr='${OP_LAN_IP_PREFIX}.20.1'
set network.iot.netmask='255.255.255.0'
commit network
EOI

config_load wireless
uci set wireless.wireless_iot_2g.network='iot'
uci commit wireless

######################################
uci batch << EOI
set firewall.zone_guest=zone
set firewall.zone_guest.name='guest'
set firewall.zone_guest.input='ACCEPT'
set firewall.zone_guest.output='ACCEPT'
set firewall.zone_guest.forward='REJECT'
set firewall.forwarding_guest=forwarding
set firewall.forwarding_guest.src='guest'
set firewall.forwarding_guest.dest='wan'
commit firewall

set dhcp.guest=dhcp
set dhcp.guest.interface='guest'
set dhcp.guest.start='100'
set dhcp.guest.limit='150'
set dhcp.guest.leasetime='12h'
commit dhcp

add_list firewall.zone_guest.network='guest'
commit firewall

set network.guest=interface
set network.guest.proto='static'
set network.guest.device='br-guest'
set network.guest.ipaddr='${OP_LAN_IP_PREFIX}.30.1'
set network.guest.netmask='255.255.255.0'
commit network
EOI

config_load wireless
uci set wireless.wireless_guest_2g.network='guest'
uci set wireless.wireless_guest_5g.network='guest'
uci commit wireless
