# ap/files/etc/uci-defaults/81-network
set -x

. /root/.env


config_load dhcp
uci del dhcp.lan.ra_slaac
uci set dhcp.lan.ignore='1'
uci commit dhcp

config_load network
uci set network.lan.ipaddr="${OP_LAN_IP_PREFIX}.1.${OP_DUMB_AP}"
uci set network.lan.gateway="${OP_LAN_IP_PREFIX}.1.1"
# uci add_list network.lan.dns="${OP_LAN_IP_PREFIX}.1.1"
uci commit network

# # not working
# uci batch << EOI
# set network.bridge_vlan_home='bridge-vlan'
# set network.bridge_vlan_home.device='br-lan'
# set network.bridge_vlan_home.vlan='1'

# set network.lan.device='br-lan.1'
# set network.lan.type='bridge'
# commit network

# set network.bridge_vlan_guest='bridge-vlan'
# set network.bridge_vlan_guest.device='br-lan'
# set network.bridge_vlan_guest.vlan='3'

# set network.guest=interface
# set network.guest.proto='static'
# set network.guest.device='br-lan.3'
# set network.guest.ipaddr='${OP_LAN_IP_PREFIX}.3.${OP_DUMB_AP}'
# set network.guest.netmask='255.255.255.0'
# set network.guest.gateway='${OP_LAN_IP_PREFIX}.3.1'
# set network.guest.type='bridge'
# commit network

# set wireless.wireless_guest_2g.network='guest'
# set wireless.wireless_guest_5g.network='guest'
# commit wireless
# EOI

uci batch << EOI
set network.switch_vlan_guest=switch_vlan
set network.switch_vlan_guest.description='guest'
set network.switch_vlan_guest.device='lan'
set network.switch_vlan_guest.vlan='3'
set network.switch_vlan_guest.vid='3'
set network.switch_vlan_guest.ports='0t 4t'
commit network

set network.br_guest=device
set network.br_guest.type='bridge'
set network.br_guest.name='br-guest'
set network.br_guest.ports='lan.3'
commit network

set network.guest=interface
set network.guest.proto='static'
set network.guest.device='br-guest'
set network.guest.netmask='255.255.255.0'
set network.guest.ipaddr="${OP_LAN_IP_PREFIX}.3.${OP_DUMB_AP}"
set network.guest.gateway="${OP_LAN_IP_PREFIX}.3.1"
commit network

set wireless.wireless_guest_2g.network='guest'
set wireless.wireless_guest_5g.network='guest'
commit wireless
EOI
