# router/files/etc/rc.local
#!/bin/sh

uci set wireless.wireless_iot.network='iot'
uci commit wireless

uci set wireless.wireless_guest.network='guest'
uci commit wireless

ifdown iot
ifdown guest
ifup iot
ifup guest

exit 0