# router/files/etc/rc.local
#!/bin/sh

uci set wireless.wireless_iot.network='iot'
uci commit wireless

ifdown iot
ifup iot

exit 0