# common/files/etc/uci-defaults/91-system
set -x

. /root/.env

if [ ! -z $OP_ROOT_PASSWORD ]; then
  (echo "$OP_ROOT_PASSWORD"; sleep 1; echo "$OP_ROOT_PASSWORD") | passwd > /dev/null
fi

# reboot every night
cat << EOF > /etc/crontabs/root
39 18 * * * sleep 130 && touch /etc/banner && reboot
EOF
chmod 0600 /etc/crontabs/root